# addons/meross_local_broker/rootfs/etc/services.d/mqtt-service-advertise/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start avahi service
# ==============================================================================
name="MQTT Local Broker"
bashio::log.info "Starting the AVAHI daemon for ${name}..."

port=$(bashio::addon.port 2001)
if [[ $port == "" ]]; then
    port=2001
fi
exec avahi-publish-service -v -f -s "$name" _meross-mqtt._tcp $port