# addons/meross_local_broker/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start nginx service
# ==============================================================================
bashio::log.info "Starting the nginx daemon"
exec nginx -c /etc/nginx/ingress.conf