# addons/meross_local_broker/rootfs/etc/services.d/meross-service/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start avahi service
# ==============================================================================
bashio::log.info "Starting the Meross daemon"
exec /usr/local/bin/meross.sh