# addons/meross_local_broker/rootfs/etc/services.d/api/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start meross API service
# ==============================================================================
bashio::log.info "Starting the Meross daemon"
exec 2>&1
exec /usr/local/bin/api.sh