# addons/meross_local_broker/rootfs/etc/services.d/avahi-public-service/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start avahi service
# ==============================================================================
name="Meross Local Broker"
bashio::log.info "Starting the AVAHI daemon for ${name}..."

exec avahi-publish-service -v -f -s "$name" _meross-lan-broker._tcp $(bashio::addon.port 2001) http_api_port=$(bashio::addon.port 2002) mqtt_port=$(bashio::addon.port 2001)