# addons/meross_local_broker/rootfs/etc/services.d/broker/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start local mqtt agent service
# ==============================================================================
bashio::log.info "Starting the mqtt client agent daemon"
exec 2>&1
exec /usr/local/bin/broker.sh