# addons/meross_local_broker/rootfs/etc/services.d/mosquitto/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start mosquitto service
# ==============================================================================
bashio::log.info "Starting the Mosquitto daemon"
exec /usr/local/bin/mosquitto.sh