# addons/meross_local_broker/rootfs/etc/services.d/avahi/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the Avahi daemon
# ==============================================================================
# shellcheck disable=SC1091

# Wait until D-Bus is available
until [[ -e /var/run/dbus/system_bus_socket ]];
do
  sleep 1
done

# Start the Avahi daemon
bashio::log.info "Starting the Avahi Daemon"
debug=$(bashio::config 'debug_mode')
if [[ $debug == true ]]; then
  exec 2>&1
  exec avahi-daemon --debug
else
  exec 2>&1
  exec avahi-daemon
fi