# addons/meross_local_broker/rootfs/etc/services.d/http-service-advertise/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start avahi service
# ==============================================================================
name="HTTP Local Broker"
bashio::log.info "Starting the AVAHI daemon for ${name}..."

port=$(bashio::addon.port 2002)
if [[ $port == "" ]]; then
    port=2002
fi
exec avahi-publish-service -v -f -s "$name" _meross-api._tcp $port